# Skripte/dts-convert/dts-convert.service
# Kopiere diese Datei nach /etc/systemd/system/
[Unit]
Description=dts-convert Service
Wants=network-online.target

[Service]
Type=oneshot
User=stef
Group=net
ExecStart=/home/$USER/github/helper/Skripte/dts-convert/dts-convert.sh

[Install]
WantedBy=multi-user.target
