# skripte/logwatch/logwatch-weekly.timer
[Unit]
Description=Wöchentlicher Log-Report (Montags 09:30)

[Timer]
OnCalendar=Mon *-*-* 09:30:00
Persistent=true

[Install]
WantedBy=timers.target
