# skripte/logwatch/logwatch-daily.timer
[Unit]
Description=Tägliche Log-Analyse (09:00)

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
